# backend/redis-worker-service.conf
[program:redisworker]
command = APP_PATH/backend/redis_worker.sh
autostart=true
autorestart=true
stdout_logfile = ROOT_DIR/var/log/gn_errors.log
redirect_stderr = true
